# script/append_txt.timer
[Unit]
Description=Executes the append_txt service every 10 sec

[Timer]
OnUnitInactiveSec=10s

[Install]
WantedBy=timers.target
